# scripts/com.brokenbats.fieldbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brokenbats.fieldbot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sam/.nvm/versions/node/v22.12.0/bin/node</string>
    <string>/Users/sam/Documents/programming_projects/broken-bats-csc/scripts/field-bot-local.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/sam/Documents/programming_projects/broken-bats-csc</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ADMIN_PASSWORD</key>
    <string>admin123</string>
    <key>PATH</key>
    <string>/Users/sam/.nvm/versions/node/v22.12.0/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Daily at 6:00 AM Eastern (11:00 UTC) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/sam/Library/Logs/field-bot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sam/Library/Logs/field-bot.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
